# evconcierge.service
# systemd unit for running the EVConcierge Django app with Gunicorn.
#
# SETUP
#   1. Replace the placeholders below:
#        <APP_USER>   -> the OS user that will run the service (e.g. ubuntu, www-data)
#        <APP_GROUP>  -> the OS group for that user
#        <APP_DIR>    -> absolute path to the Django project (the folder containing manage.py)
#      <APP_DIR> must contain a `.venv/` virtualenv with gunicorn installed
#      (created via `uv sync` from the project root).
#
#   2. Install the unit and enable it:
#        sudo cp <APP_DIR>/evconcierge.service /etc/systemd/system/evconcierge.service
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now evconcierge
#
#   3. Inspect:
#        sudo systemctl status evconcierge
#        sudo journalctl -u evconcierge -f
#
# NOTES
#   - Gunicorn binds to a Unix socket at /run/evconcierge/evconcierge.sock.
#     RuntimeDirectory= makes systemd create /run/evconcierge/ on start (owned
#     by <APP_USER>:<APP_GROUP>, mode 0750) and remove it on stop.
#   - Point your reverse proxy at the socket. Example nginx upstream:
#       upstream evconcierge { server unix:/run/evconcierge/evconcierge.sock; }
#     Make sure the nginx user (e.g. www-data) is in <APP_GROUP> so it can
#     read the socket, or relax RuntimeDirectoryMode (e.g. 0755).
#   - Adjust --workers based on CPU cores (a common rule: 2 * cores + 1).
#   - Logs go to journald via stdout/stderr (--access-logfile - / --error-logfile -).

[Unit]
Description=EVConcierge Django app (Gunicorn)
After=network.target

[Service]
Type=simple
User=<APP_USER>
Group=<APP_GROUP>
WorkingDirectory=<APP_DIR>
Environment="PATH=<APP_DIR>/.venv/bin"
RuntimeDirectory=evconcierge
RuntimeDirectoryMode=0750
ExecStart=<APP_DIR>/.venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/run/evconcierge/evconcierge.sock \
    --access-logfile - \
    --error-logfile - \
    EVConcierge.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
